# deploy_tools/gunicorn-upstart.template.conf
dcription "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid elspeth

chdir /home/elspeth/ssuperlists-staging.ottg.eu/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/superlists-staging.ottg.eu.socket superlists.wsgi:application
